# systemd/unbound.service.d/override.conf
[Service]
ExecStartPost=/bin/sleep 2
ExecStartPost=/usr/local/bin/unbound-cache-load
ExecStop=/usr/local/bin/unbound-cache-dump
