# systemd/unbound.service.d/override.conf
# Deploys to: /etc/systemd/system/unbound.service.d/override.conf
# Adds cache load on start (2s delay for socket) and cache dump on stop.

[Service]
ExecStartPost=/bin/sleep 2
ExecStartPost=/usr/local/bin/unbound-cache-load
ExecStop=/usr/local/bin/unbound-cache-dump
