# systemd/cake.service
# Deploys to: /etc/systemd/system/cake.service
# Reload: sudo systemctl daemon-reload && sudo systemctl enable --now cake

[Unit]
Description=CAKE queue discipline (upload bufferbloat management)
After=network-online.target wg-quick@wg0.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cake-setup.sh
ExecStop=/sbin/tc qdisc del dev enp1s0 root
ExecStop=/bin/bash -c 'for proto in udp tcp; do for sport in 53 5335; do iptables -t mangle -D POSTROUTING -p $proto --sport $sport -j DSCP --set-dscp-class CS5 2>/dev/null || true; done; done'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
