# systemd/cloudflared-proxy-dns.service
# Deploys to: /etc/systemd/system/cloudflared-proxy-dns.service
# Reload: sudo systemctl daemon-reload && sudo systemctl enable --now cloudflared-proxy-dns
#
# THE CLOUDFLARE ENDOSYMBIONT. A contained cloudflared organelle that lives on
# the t630's loopback and does exactly ONE job: speak DNS-over-HTTPS to
# Cloudflare's anycast network on behalf of Unbound's STREAMING forward-path.
#
# Why this is endosymbiosis (an engulfed organelle), not a takeover:
#   * Listens only on 127.0.0.1:5053 — nothing off-box can reach it.
#   * ONLY the forward-zones in unbound/streaming-forward.conf point at it. The
#     recursive path — banking, email, health, the default — NEVER touches it
#     and stays sovereign: Cloudflare gets zero visibility into sensitive lookups.
#     That invariant is the whole point; do not widen what feeds this proxy.
#   * Upstreams are IP-LITERAL DoH endpoints, so the organelle needs no DNS of
#     its own to bootstrap (avoids the chicken-and-egg where the thing that
#     resolves names would first have to resolve a name). Cloudflare's TLS cert
#     carries the 1.1.1.1 / 1.0.0.1 IP SANs, so validation succeeds against the IP.
#
# What it buys the forward-path: the ISP sees only encrypted HTTPS on :443
# instead of plaintext UDP/53 streaming lookups, and queries ride Cloudflare's
# anycast edge. The privacy-for-speed trade was already accepted for these
# domains — this just encrypts the hop the ISP previously read in the clear.
#
# If the organelle dies, Unbound's `forward-first: yes` falls back to recursive
# resolution for the streaming domains. Slower, but never a hard failure.

[Unit]
Description=cloudflared DoH proxy (Cloudflare endosymbiont for the streaming forward-path)
After=network-online.target
Wants=network-online.target
Before=unbound.service

[Service]
ExecStart=/usr/local/bin/cloudflared proxy-dns \
    --address 127.0.0.1 \
    --port 5053 \
    --upstream https://1.1.1.1/dns-query \
    --upstream https://1.0.0.1/dns-query
Restart=on-failure
RestartSec=5
# Containment hardening — the organelle reads/writes nothing and runs unprivileged.
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
