# 11-console/ttyd-thinclient.service
# Deploys to: /etc/systemd/system/ttyd-thinclient.service
# Reload:     sudo systemctl daemon-reload && sudo systemctl enable --now ttyd-thinclient.service
#
# A browser login shell on the t630 itself (the thin client), on :7681, gated to LAN + WG by UFW.
# ttyd's --credential is HTTP basic auth and is the ONLY thing between the
# network and a shell as USER — treat it like a root password (set it in /etc/a777ance/ttyd.env).
# Hardening ladder is in 11-console/README.md (TLS with -S, OS `login` instead of bash, etc.).
#
#   sudo apt install ttyd          # Ubuntu 24.04 universe
#   replace USER below with the t630 login user
[Unit]
Description=A777ance console — web terminal to the thin client (t630) via ttyd
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=USER
EnvironmentFile=/etc/a777ance/ttyd.env
# --writable lets you type into the shell; --max-clients caps concurrent sessions.
ExecStart=/usr/bin/ttyd --port 7681 --interface 0.0.0.0 --credential ${TTYD_CRED} --max-clients 3 --writable bash -l
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
