# 11-console/ttyd-laptop.service
# Deploys to: /etc/systemd/system/ttyd-laptop.service
# Reload:     sudo systemctl daemon-reload && sudo systemctl enable --now ttyd-laptop.service
#
# A browser login shell on the LINUX LAPTOP, served from the t630 on :7682 and gated to LAN + WG.
# The t630 is the jump host: ttyd runs `ssh ${LAPTOP_SSH}`, so the page is reachable whenever the
# t630 is up, and it connects through to the laptop when the laptop is reachable (LAN or WG).
#
# SSH auth is interactive by default — the page prompts for the laptop's password, so no key is
# stored on the t630. To go passwordless instead, drop an SSH key in USER's ~/.ssh on the t630
# and authorise it on the laptop (documented in 11-console/README.md).
#
# Set LAPTOP_SSH (user@host) in /etc/a777ance/ttyd.env. Prefer a STABLE address — the laptop's
# WireGuard IP (10.8.0.x) or a DHCP-reserved LAN IP, not a floating lease. Replace USER below.
[Unit]
Description=A777ance console — web terminal to the Linux laptop (SSH jump from the t630) via ttyd
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=USER
EnvironmentFile=/etc/a777ance/ttyd.env
ExecStart=/usr/bin/ttyd --port 7682 --interface 0.0.0.0 --credential ${TTYD_CRED} --max-clients 2 --writable ssh ${LAPTOP_SSH}
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
