# 01-core-network/unbound/unbound-cache-dump.service
# Deploys to: /etc/systemd/system/unbound-cache-dump.service
# Triggered by: unbound-cache-dump.timer (hourly) and unbound.service stop

[Unit]
Description=Dump Unbound cache
Requires=unbound.service
After=unbound.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/unbound-cache-dump
